# cbv.socket
[Unit]
Description=CBV socket

[Socket]
ListenStream=23