# cbv.socket
[Unit]
Description=CBV socket

[Socket]
ListenStream=23

[Install]
WantedBy=sockets.target